# scripts/com.tulsbot.context-report.evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tulsbot.context-report.evening</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tulioferro/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/tulioferro/.openclaw/workspace/scripts/context-report.ts</string>
        <string>evening</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/tulsbot-report-evening.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tulsbot-report-evening.error</string>
</dict>
</plist>