# scripts/com.openclaw.heartbeat-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.heartbeat-loop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>--import</string>
        <string>tsx</string>
        <string>scripts/run-heartbeat.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tulioferro/.openclaw/workspace</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tulioferro/.openclaw/logs/heartbeat-loop-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tulioferro/.openclaw/logs/heartbeat-loop-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tulioferro/Library/pnpm:/Users/tulioferro/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
