# scripts/com.openclaw.heartbeat-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.heartbeat-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pnpm</string>
        <string>tsx</string>
        <string>/Users/tulioferro/.openclaw/workspace/scripts/run-heartbeat-daily.ts</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/tulioferro/.openclaw/workspace</string>
    <key>StandardOutPath</key>
    <string>/Users/tulioferro/.openclaw/logs/heartbeat-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tulioferro/.openclaw/logs/heartbeat-daily.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
